# wakeupbot.service
[Unit]
Description=WakeupBot

[Service]
Type=simple
ExecStart=/bin/node /home/mjw/groupme/srv.js

[Install]
WantedBy=multi-user.target
